# deploy/labubu.service
[Unit]
Description=Labubu Animation FastAPI (uvicorn)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/srv/labubu-animation/backend
EnvironmentFile=/srv/labubu-animation/backend/.env
ExecStart=/usr/bin/env uvicorn app:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
